# server/drawpile-srv.service
[Unit]
Description=Drawpile dedicated server
After=network.target

[Service]
ExecStart=/usr/bin/drawpile-srv
Type=notify
NotifyAccess=main
User=drawpile

[Install]
WantedBy=multi-user.target

