# desktop/drawpile-srv.socket
[Unit]
Description=Drawpile dedicated server socket
After=network.target

[Socket]
ListenStream=27750

[Install]
WantedBy=sockets.target

